# ops/server_templates/copybot-discovery-v2-publish.service
[Unit]
Description=Copybot Discovery V2 Publish
After=solana-copy-bot.service

[Service]
Type=oneshot
User=copybot
WorkingDirectory=/var/www/solana-copy-bot
EnvironmentFile=/etc/solana-copy-bot/app.env
ExecStart=/usr/bin/flock /var/lock/copybot-discovery-v2.lock /var/www/solana-copy-bot/bin/discovery_v2_publish --config /etc/solana-copy-bot/live.server.toml --commit --acknowledge-daemon-restart-required
TimeoutStartSec=1200
