# ops/server_templates/copybot-discovery-v2-watchdog.timer
[Unit]
Description=Run Copybot Discovery V2 Watchdog every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=copybot-discovery-v2-watchdog.service

[Install]
WantedBy=timers.target
