# ops/server_templates/copybot-discovery-recent-raw-snapshot.timer
[Unit]
Description=Copybot Discovery Recent Raw Journal Snapshot Timer

[Timer]
OnBootSec=7m
# The snapshot service enforces its own cadence and can now exit with a
# transient success (`75`) for `deferred` / `retryable_busy`, so keeping the
# timer tighter than the steady-state snapshot cadence gives bounded retries
# without operator intervention once the service has exited via its own
# attempt-duration budget or contention contract.
OnUnitActiveSec=5m
RandomizedDelaySec=30s
Persistent=true
Unit=copybot-discovery-recent-raw-snapshot.service

[Install]
WantedBy=timers.target
