# ops/server_templates/copybot-discovery-v2-publish.timer
[Unit]
Description=Copybot Discovery V2 Publish Timer

[Timer]
OnBootSec=45m
OnUnitActiveSec=60m
RandomizedDelaySec=5m
Persistent=true
Unit=copybot-discovery-v2-publish.service

[Install]
WantedBy=timers.target
