# ops/server_templates/copybot-discovery-v2-prepare-quality.timer
[Unit]
Description=Copybot Discovery V2 Quality Prepare Timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=30m
RandomizedDelaySec=2m
Persistent=true
Unit=copybot-discovery-v2-prepare-quality.service

[Install]
WantedBy=timers.target
