# ops/server_templates/copybot-discovery-recent-raw-snapshot.service
[Unit]
Description=Copybot Discovery Recent Raw Journal Snapshot
After=solana-copy-bot.service

[Service]
Type=oneshot
User=copybot
WorkingDirectory=/var/www/solana-copy-bot
EnvironmentFile=/etc/solana-copy-bot/app.env
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
SuccessExitStatus=75
SyslogIdentifier=discovery_recent_raw_snapshot
StandardOutput=journal
StandardError=journal
TimeoutStartSec=10min
ExecStart=/var/www/solana-copy-bot/bin/discovery_recent_raw_snapshot --config /etc/solana-copy-bot/live.server.toml --scheduled --json
