# ops/server_templates/copybot-discovery-v2-prepare-quality.service
[Unit]
Description=Copybot Discovery V2 Quality Prepare
After=solana-copy-bot.service

[Service]
Type=oneshot
User=copybot
WorkingDirectory=/var/www/solana-copy-bot
EnvironmentFile=/etc/solana-copy-bot/app.env
ExecStart=/usr/bin/flock -n -E 75 /var/lock/copybot-discovery-v2.lock /var/www/solana-copy-bot/bin/discovery_v2_prepare_quality --config /etc/solana-copy-bot/live.server.toml --max-mints 20000 --commit --materialize-status
SuccessExitStatus=75
TimeoutStartSec=900
